# .config/systemd/user/x11vnc.service
[Unit]
Description=x11vnc VNC server for the current X session
After=graphical-session.target

[Service]
Type=simple
Environment=DISPLAY=:1
Environment=XAUTHORITY=/run/user/1000/gdm/Xauthority
ExecStart=/usr/bin/x11vnc -display :1 -auth /run/user/1000/gdm/Xauthority -clip 2560x1440+0+0 -forever -loop -noxdamage -repeat -passwd test12 -listen 0.0.0.0 -rfbport 5900 -shared -wait 5 -defer 5 -speeds lan -threads
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
