# Library/LaunchAgents/com.simonrw.theming.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.simonrw.theming</string>
	<key>ProgramArguments</key>
	<array>
        <string>/opt/homebrew/opt/python/bin/python3</string>
        <string>/Users/simon/dotfiles/theming/theming.py</string>
	</array>
	<key>WatchPaths</key>
	<array>
		<string>/Users/simon/Library/Preferences/.GlobalPreferences.plist</string>
	</array>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
